# rpi/settings/crontab
0 * * * * cd /home/pi/diy-baby-crib-mobile/rpi/facedetection; python facedetection.py
